# configs/monthly-compression.cron.template
# Compress old data on 1st of month at 2 AM
0 2 1 * * ${DEPLOY_DIR}/scripts/compress-old-data.sh >> ${DEPLOY_DIR}/logs/compression.log 2>&1

